# sysadmin/crai-cron.service
[Unit]
Description=crai-cron
After=network.target

[Service]
ExecStart=@crai-cron@ --database=@database@ --mirror=@mirror@
Type=oneshot
User=crai
